# board/opendingux/gcw0/target_skeleton/etc/network/if-up.d/udhcpd
#!/bin/sh

[ "$IFACE" == "usb0" ] || exit 0

/sbin/start-stop-daemon -S -q -b -m -p /var/run/udhcpd.pid \
	-x /usr/sbin/udhcpd -- -f -S /etc/udhcpd.conf
